# service/osfooler-ng.conf.upstart
#!upstart
#
description "OSfooler upstart script"
 
start on startup
stop on shutdown
 
script
    /usr/bin/python3 OSFooler-ng/osfooler-ng.py -m "${__OS}" -o "${__OSGENRE}" -d "${__DETAILS_P0F}"
end script
 
pre-stop script
    # Here we're using the pre-stop script to stop the application
    exec forever stop $APPLICATION_PATH
end script
